# templates/redis-server.supervisor.conf.erb
[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
stderr_logfile=/var/log/redis/error.log
autostart=true
autorestart=true
stopwaitsecs=10
user=redis
group=redis
